# tools/99-avm.rules
SUBSYSTEM=="block", ACTION=="add" RUN+="/home/pi/mdesk/tools/mount.sh"
